# telegram-tracker.service
[Unit]
Description=Telegram Online Status Tracker
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/telegram-tracker
ExecStart=/usr/bin/python3 /opt/telegram-tracker/telegram_tracker.py
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=telegram-tracker

[Install]
WantedBy=multi-user.target